# ansible/roles/ps4-controller/templates/ds4drv.service.j2
[Unit]
Description=ds4drv daemon
Requires=bluetooth.service
After=bluetooth.service

[Service]
EnvironmentFile=/etc/donkeycar/ds4drv.env
ExecStart=sudo /home/pi/venv/bin/ds4drv --controller-limit 1 --led ${controller_color}
Restart=on-abort

[Install]
WantedBy=bluetooth.target