# ansible/roles/shutdown-btn/templates/button-shutdown.service.j2
[Unit]
Description=button shutdown service

[Service]
ExecStart={{ shutdown_btn_script_path }}
Restart=on-failure

[Install]
WantedBy=multi-user.target
